# bootc-1.16.2/systemd/bootc-destructive-cleanup.service
[Unit]
Description=Cleanup previous the installation after an alongside installation
Documentation=man:bootc(8)

[Service]
Type=oneshot
ExecStart=/usr/lib/bootc/fedora-bootc-destructive-cleanup
PrivateMounts=true

# No [Install] section, this is enabled via generator
