# bootc-1.16.2/systemd/bootc-status-updated-onboot.target
[Unit]
Description=Bootc status trigger state sync
Documentation=man:bootc-status-updated.target(8)
After=sysinit.target

# No [Install] section, this is enabled via generator
